# pi/deploy/matau-gps-watchdog.service
[Unit]
Description=Matau GPS watchdog — revive stalled PL2303 USB GPS
After=signalk.service

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /usr/local/bin/matau_gps_watchdog.py
